# etc/inittab
::sysinit:/etc/rc.d/rcS
ttyS0::respawn:/bin/sh
